# systemd/eclipanel-backend.service
#! /bin/sh
[Unit]
Description=EcliPanel Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/samba/shared/EcliPanel/v3/backend
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -lc './start.sh'
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
KillMode=control-group
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eclipanel-backend

[Install]
WantedBy=multi-user.target
