# systemd/eclipanel-frontend.service
[Unit]
Description=EcliPanel Frontend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/samba/shared/EcliPanel/v3/frontend
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=5
Environment=NODE_ENV=development
Environment=BACKEND_URL=http://localhost:3001
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
